# pkg/debian/DEBIAN/prerm
#!/bin/bash

PATH=/sbin:/bin:/usr/sbin:/usr/bin

case "$1" in
  *)
    systemctl disable --now flecs-webapp
    ;;
esac

exit 0
