# debian/DEBIAN/postinst
#!/bin/bash

PATH=/sbin:/bin:/usr/sbin:/usr/bin

case "$1" in
  *)
    docker load --input /opt/flecs-webapp/assets/flecs-webapp_*.tar.gz
    systemctl enable --now flecs-webapp
    ;;
esac

exit 0
