# debian/DEBIAN/postinst
#!/bin/bash

PATH=/sbin:/bin:/usr/sbin:/usr/bin

case "$1" in
  *)
    systemctl enable --now flecs-webapp
    ;;
esac

exit 0
